# upgrade-files/rootfs/etc/init.d/agingtest.h3
#!/bin/sh /etc/rc.common

START=90

start() {
   /usr/bin/aging_test 4 131072 32 210 
}
