# upgrade-files/rootfs/etc/crontabs/root
*/15   *  *   *   *     /usr/bin/cgminer-monitor
0    */6  *   *   *     has_data_bak=`mount | grep data_bak` && if [ "$has_data_bak" != "" ]; then rsync -ar --delete /data/etc /data_bak; fi
0    */6  *   *   *     suffix=`tail -c 3 /tmp/miner-info` && if [[ -f /tmp/*.core && $suffix != "CD" ]]; then echo ".CD" >> /tmp/miner-info; fi
